# telegram-reboot-alert.service
[Unit]
Description=Send Telegram boot and shutdown alerts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 15
ExecStart=/opt/telegram-reboot-alert/send-alert.sh boot
ExecStop=/opt/telegram-reboot-alert/send-alert.sh shutdown
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
